# wp12-support/ansible/slack_sync/slack-sync.service
[Unit]
Description=Slack Syncronization Service
After=syslog.target network.target 

[Service]
Type=simple
Restart=always
User=app
Group=app
ExecStart=python3 /home/app/slack-sync.py
EnvironmentFile=/etc/environment

[Install]
WantedBy=multi-user.target